# db-backup/crontab
0 3 * * * sh /usr/src/app/backupdb.sh >> /var/log/cron.log 2>&1
